# hotspot_monitor.service
[Unit]
Description=Hotspot Monitor Service

[Service]
Type=simple
ExecStart=/home/pi/monitor/monitor.sh

[Install]
WantedBy=multi-user.target
